# zaplinkcore.service
[Unit]
Description=ZapLinkCore - OTA TV Streaming Server
After=network.target

[Service]
Type=simple
User=zaplink
Group=zaplink
WorkingDirectory=/opt/zaplink
ExecStart=/opt/zaplink/zaplinkcore
Restart=on-failure
RestartSec=5

# Allow access to DVB adapters
SupplementaryGroups=video

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/zaplink

[Install]
WantedBy=multi-user.target
